# kurento-room-demo/package/support-files/kroomdemo.service
[Unit]
Description=KRoomDemo daemon
After=network.target

[Service]
User=nobody
ExecStart=/usr/bin/kroomdemo

[Install]
WantedBy=default.target